# virtual-device/conf/supervisord.conf
[supervisord]
nodaemon=true

[program:app]
directory=/app
command=python app.py
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
autorestart=true